# wg-vpn-auto-restore.service
[Unit]
Description=Restore WireGuard VPN Auto on boot
After=NetworkManager-wait-online.service
Requires=NetworkManager-wait-online.service

[Service]
Type=oneshot
ConditionPathExists=/var/lib/wg-vpn-auto/enabled
ExecStart=/bin/bash -c '. /etc/wg-vpn-auto.conf && exec nmcli connection up "$DUMMY_CONN_ID"'

[Install]
WantedBy=multi-user.target
