# wg-vpn-auto.service
[Unit]
Description=WireGuard VPN Auto Orchestrator
After=NetworkManager.service
Requires=NetworkManager.service

[Service]
Type=simple
ExecStart=/usr/local/bin/wg-vpn-auto.sh --daemon
Restart=on-abnormal
RestartSec=5
